# etc/systemd/system/proxmox-auto-update.service
[Unit]
Description=Proxmox Automatic System Update and Reboot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/proxmox-auto-update.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
